# docker/supervisord.ini
[supervisord]
nodaemon=true

[program:sshd]
autorestart=true
command=/usr/sbin/sshd -D
directory=%(HOME)s
stdout_logfile=syslog
stderr_logfile=syslog
user=root

[program:rsyslog]
autorestart=true
command=rsyslogd -n
directory=%(HOME)s
stdout_logfile=syslog
stderr_logfile=syslog
user=root
